# archive/shelved-2026-03-21/things-bridge/com.selene.things-bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.selene.things-bridge</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/chaseeasterling/selene/scripts/things-bridge/process-pending-tasks.sh</string>
    </array>

    <key>WatchPaths</key>
    <array>
        <string>/Users/chaseeasterling/selene/vault/things-pending</string>
    </array>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/chaseeasterling/selene/logs/things-bridge-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/chaseeasterling/selene/logs/things-bridge-stderr.log</string>
</dict>
</plist>
